# deploy/polyasr.service.template
[Unit]
Description=polyasr server (CUDA / Qwen3-ASR + forced aligner)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__REPO__/cuda
Environment=POLYASR_MODEL=Qwen/Qwen3-ASR-1.7B
Environment=POLYASR_ALIGNER_MODEL=Qwen/Qwen3-ForcedAligner-0.6B
Environment=POLYASR_DEVICE=cuda:0
Environment=POLYASR_DTYPE=bfloat16
Environment=POLYASR_PORT=8766
Environment=POLYASR_LOG_DIR=__REPO__/cuda/logs
# Idle-evict the resident model after this many seconds so a co-resident
# co-resident workload can reclaim the GPU. Default 0 = never evict: this
# service backs benchday live dictation, which needs the model warm at all times.
Environment=POLYASR_IDLE_EVICT_SECONDS=0
# Hugging Face cache persists across restarts so model reload is fast.
Environment=HF_HOME=__HF_HOME__
# huggingface.co is unreachable from CN hosts (SSL EOF). Even via the mirror,
# the model's speculative `custom_generate/generate.py` Hub fetch hangs on the
# CDN and wedges model load (weights load fine, then from_pretrained hangs).
# The model+config are cached, so load fully offline; the mirror endpoint is a
# fallback for the rare case the cache must be repopulated manually.
Environment=HF_ENDPOINT=https://hf-mirror.com
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1
ExecStart=__REPO__/cuda/venv/bin/python __REPO__/cuda/server.py
Restart=on-failure
RestartSec=5
StandardOutput=append:__LOG_DIR__/polyasr.log
StandardError=append:__LOG_DIR__/polyasr.err

[Install]
WantedBy=multi-user.target
