# launchd/io.zigzag.polyasr.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.zigzag.polyasr</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VENV__/bin/python</string>
        <string>__REPO__/server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <!-- Model: Qwen/Qwen3-ASR-0.6B (fast) or Qwen/Qwen3-ASR-1.7B (best) -->
        <key>POLYASR_MODEL</key>
        <string>Qwen/Qwen3-ASR-1.7B</string>
        <!-- Idle-evict the resident model after N seconds so a co-resident
             polytts/renderer can reclaim memory. 0 disables (never evict). -->
        <key>POLYASR_IDLE_EVICT_SECONDS</key>
        <string>180</string>
        <!-- Use HF mirror when huggingface.co is unreachable -->
        <key>HF_ENDPOINT</key>
        <string>https://hf-mirror.com</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__REPO__/server.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__REPO__/server.stderr.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
